# cron.daily/letsencrypt
#!/bin/sh

exec /usr/bin/letsencrypt renew
